# automatic.cron
# crontab for web-sysinfo user
  * *  *   *   *  echo "hello" | wall

